# launchagents/Library/LaunchAgents/com.lucas.bd-evening.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!--
  TEMPLATE — rendered to real plist by setup.sh's render_launchagents().
  __USER__ is substituted with the current $USER at install time.
  The rendered file lands at ~/Library/LaunchAgents/<filename>.plist.

  Time-of-day BetterDisplay mode scheduler. Mode parameter is the second
  argument to bd-apply.sh (dawn/day/afternoon/evening/night).
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.lucas.bd-evening</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__DOTFILES_DIR__/scripts/scripts/bd-apply.sh</string>
        <string>evening</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>18</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>RunAtLoad</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/tmp/bd-launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/bd-launchd.log</string>
</dict>
</plist>
