# launchagents/Library/LaunchAgents/com.lucas.wallpaper-rotate.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!--
  TEMPLATE — rendered to real plist by setup.sh's render_launchagents().
  __USER__ is substituted with the current $USER at install time.
  The rendered file lands at ~/Library/LaunchAgents/<filename>.plist.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.lucas.wallpaper-rotate</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-lc</string>
        <string>/Users/__USER__/scripts/wallpaper-rotate.sh</string>
    </array>

    <!-- Fire every 3600s (1 hour). -->
    <key>StartInterval</key>
    <integer>3600</integer>

    <!-- Also fire once at load so the wallpaper updates immediately on login. -->
    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/__USER__/Library/Logs/wallpaper-rotate.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/__USER__/Library/Logs/wallpaper-rotate.log</string>

    <!-- Keep the launchd job alive across sleep/wake. -->
    <key>KeepAlive</key>
    <false/>

    <!-- Reasonable nice value — wallpaper changes are not urgent. -->
    <key>Nice</key>
    <integer>10</integer>
</dict>
</plist>
